# jmxtrans/templates/jmxtrans.init.erb
#!/bin/sh
# chkconfig: 345 90 10
# description: jmxtrans Startup Script
# processname: java

start() {
  <%= @jmxtrans_home -%>/bin/start_jmxtrans.sh kafka-server
}

stop() {
  <%= @jmxtrans_hoem -%>/bin/stop_jmxtrans.sh
}

usage() {
  echo "This script supports {start,stop}"
}

case "$1" in

  start)
    start
    ;;
  stop)
    stop
    ;;
  *)
    usage
    ;;
esac
