# new_base/modules/base/templates/profile.d/titlebar.sh.erb
#!/bin/bash
#(@) Joe DeCello, Symantec Cloud Platform Engineering
#######################
## Managed by Puppet ##
#######################
TITLE="[<%= @datacenter %>] $USER@<%= @fqdn %> [<%= @ipaddress %>]"
PROMPT_COMMAND='echo -ne "\033]0; $TITLE \007"'
